# scripts/systemd/mcp-backend.service
[Unit]
Description=MCP Jira Thing backend API
After=network.target

[Service]
Type=simple
User=grimlock
Group=grimlock
WorkingDirectory=/var/www/vhosts/mcp-jira-thing.truvis.co
EnvironmentFile=/etc/mcp-jira-thing/config.ini
ExecStart=/var/www/vhosts/mcp-jira-thing.truvis.co/mcp-backend
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target


